# docker/supervisord.railway.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

; ── FastAPI (uvicorn) — internal port 8000 ──────────────────────────────────
[program:uvicorn]
command=uvicorn app.main:app
    --host 127.0.0.1
    --port 8000
    --workers 1
    --loop uvloop
directory=/app/backend
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ── Voice Worker — background job processor ─────────────────────────────────
[program:voice-worker]
command=python -m app.workers.voice_worker
directory=/app/backend
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

; ── nginx — public port 80 (starts after uvicorn is up) ─────────────────────
[program:nginx]
command=nginx -c /etc/nginx/nginx.conf -g "daemon off;"
autostart=true
autorestart=true
startretries=3
; Give uvicorn a moment to start before nginx accepts requests
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
